# ec-server-update.service
[Unit]
Description=update and restart ec server

[Service]
Type=simple
ExecStart=/mnt/EC_server/run.sh